# ansible/files/superapp.service
[Unit]
Description=SuperApp

[Service]
ExecStart=/home/ubuntu/bin/reddit

[Install]
WantedBy=multi-user.target